# 80-setacl-helper.rules
ACTION=="add|change", ENV{DEVNAME}!="", RUN+="/usr/lib/setacl-helper/script $env{DEVNAME}"
